# deploy/seiyuumatch@.service
[Unit]
Description=SeiyuuMatch recognition worker on port %i
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/SeiyuuMatch
ExecStart=/root/miniconda3/envs/seiyumatch/bin/python -u server.py --host 127.0.0.1 --port %i
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
Environment=OMP_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
Environment=RECOGNITION_QUEUE_TIMEOUT=30
Environment=JOB_QUEUE_MAX_SIZE=200
Environment=MAX_IMAGE_DIM=1280

[Install]
WantedBy=multi-user.target
